# nova-consoleauth.conf
[program:nova-api]
directory=/
command=service nova-consoleauth start
user=root
autostart=true
autorestart=true
priority=999
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log